# scripts/sse100.crontab
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.XXXXBrkP2p installed on Mon Aug 12 17:29:11 2013)
# (Cron version V5.0 -- $Id: crontab.c,v 1.12 2004/01/23 18:56:42 vixie Exp $)
#JR - restart TomCat every hour. There is a java thread buildup problem
36 * * * * ${HOME}/scripts/restart_tomcat
02 17 * * * ${HOME}/sonata_install/bin/obs-summary-report.sh -dbhost sse100 -dbname kepler8ghz -hours 24 | mailx -s "SonATA daily obs report - kepler8ghz" -r observing@seti.org observing@seti.org 


06 17 * * * ${HOME}/sonata_install/bin/obs-summary-report.sh -dbhost sse100 -dbname galcenter2 -hours 24 | mailx -s "SonATA daily obs report - galcenter2" -r observing@seti.org observing@seti.org 
08 17 * * * ${HOME}/sonata_install/bin/galcenter2-keplerHz-progress.sh 6970 7045| mailx -s "SonATA daily galcenter2 progress report" -r observing@seti.org observing@seti.org 
08 17 * * * ${HOME}/sonata_install/bin/obs-cal-report.sh | mailx -s "SonATA calibration report" -r observing@seti.org observing@seti.org
# 10 14 * * * ${HOME}/sonata_install/bin/hat-creek-ntp-server-ping-test.sh | mailx -s "Hat Creek SonATA NTP Server ping test" -r observing@seti.org observing@seti.org
12 17 * * * ${HOME}/sonata_install/bin/daily-error-summary.sh | mailx -s "SonATA Error Summary" -r observing@seti.org observing@seti.org
01 00 * * * ${HOME}/sonata_install/bin/sse100-backup.sh | mailx -s "SonATA Daily Backup" -r observing@seti.org observing@seti.org
#30 00 * * * ${HOME}/sonata_install/bin/mybackup.sh | mailx -s "SonATA Daily Backup" -r observing@seti.org observing@seti.org
12 17 * * * ${HOME}/sonata_install/bin/db-backup-sse100.sh | mailx -s "SonATA Database Backup" -r observing@seti.org observing@seti.org
# Every minute, start the seekerWatcher if it is not running - JR - Jan 04, 2012
* * * * * cd /home/sonata/jrichards/dev/seekerWatch; /home/sonata/jrichards/dev/seekerWatch/seekerStart
# Every minute, start Charlie if it is not running - JR - July 03, 2012
* * * * * cd /home/sonata/jrichards/dev/seekerWatch; /home/sonata/jrichards/dev/seekerWatch/charlieStart
# Every minute, start TwitterManager if it is not running - JR - Oct 27, 2012
* * * * * cd /home/sonata/jrichards/dev/seekerWatch; /home/sonata/jrichards/dev/seekerWatch/twitterStart
# Every 24 hours run the obs report that sends to setiquest.info
00 13 * * * cd /home/sonata/jrichards/dev/seekerWatch; /home/sonata/jrichards/dev/seekerWatch/dailyObsReportStart

#

#

#schdule Anu 08/13 08/14 
# starting with galcenter 
#05 02 14 08 * /home/sonata/sonata_install/scripts/sonata-cron-begin-galcenter-session.sh
#commensal
#06 02 14 08  * /home/sonata/scripts/start_commensal 6.90 7082 7082 
#GALCENTER -- KEPLER
#22 05 14 08 * /home/sonata/sonata_install/scripts/sonata-cron-switch-to-kepler.sh
#KEPLER -- SRI
#55 13 14 08 * /home/sonata/sonata_install/scripts/sonata-cron-end-obs-session.sh
